# utils/fromchat.service
[Unit]
Description=FromChat server
After=multi-user.target
Wants=network-online.target
After=network-online.target
StartLimitIntervalSec=60
StartLimitBurst=3

[Service]
Type=simple
User=root
Group=root
ExecStart=/bin/bash -c "COMPOSE_PROFILES=production docker compose up --remove-orphans --force-recreate"
ExecStop=/bin/bash -c "COMPOSE_PROFILES=production docker compose down --remove-orphans"
WorkingDirectory=/home/denis0001-dev/actions-runner/_work/FromChat/FromChat
Restart=always
RestartSec=10
StartLimitBurst=3

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/var/log

# Logging
StandardOutput=journal
StandardError=journal
StandardInput=tty-force
SyslogIdentifier=fromchat

[Install]
WantedBy=multi-user.target